# live-server/srv1782156/systemd/goat-intel.service
[Unit]
Description=GOAT Intel - Money Penny and crew brain
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/goat/intel
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 goat_intel.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
